# ansible/templates/navu-house-backup-trigger.service.j2
[Unit]
Description=Run Navu House backup on demand

[Service]
Type=oneshot
ExecStart=/usr/local/bin/navu-house-backup-trigger.sh
